# mpd_oled.service
[Unit]
Description=MPD OLED Display

[Service]
ExecStart=/usr/local/bin/mpd_oled -o 6 -b 21 -g 1 -f 15

[Install]
WantedBy=multi-user.target
